# deploy/systemd/compania-backend.service.example
[Unit]
Description=Compania Service Backend
After=network.target

[Service]
Type=simple
User=compania
WorkingDirectory=/opt/compania_service/backend
EnvironmentFile=/opt/compania_service/backend/.env
ExecStart=/usr/bin/node dist/src/main.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
